# root/etc/services.d/deemix/run
#!/usr/bin/with-contenv bash

UMASK_SET=${UMASK_SET:-022}
umask "$UMASK_SET"

cd /deemix-gui/server

HOME=/deem

if [ ! -z ${INTPORT} ]; then
	port=$INTPORT
else
	port=6595
fi


echo "[services.d] Starting Deemix"
s6-setuidgid abc yarn start-build --host 0.0.0.0 --port $port